# webapp/deploy/thermomix-webapp.service
[Unit]
Description=thermomix-master Next.js webapp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/thermomix-master/webapp
Environment=NODE_ENV=production
Environment=THERMOMIX_STATE_DIR=/var/lib/thermomix-master
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3000
Restart=always
RestartSec=5
# SSE-Streams (chat /chat, offene Browser-Tabs) halten Next.js fest → SIGTERM hängt.
# Nach 10s SIGKILL erzwingen, statt ~60s im stop-sigterm zu sitzen (nginx 502-Fenster).
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
