# webapp/deploy/thermomix-autoupdate.timer
[Unit]
Description=thermomix-master auto-update timer (every 10 min)

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
