# webapp/deploy/thermomix-chat.service
[Unit]
Description=thermomix-master embedded Claude-Code chat (ttyd + tmux + claude)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/thermomix-master
Environment=HOME=/root
Environment=TERM=xterm-256color
# ttyd listens on localhost:7681; nginx reverse-proxies /chat → :7681.
# -W = writable (user can type), --max-clients 4 = small abuse cap.
# tmux session "thermomix" persists across browser refreshes.
ExecStart=/usr/local/bin/ttyd -W --max-clients 4 -p 7681 -i 127.0.0.1 \
  /usr/bin/tmux new -A -s thermomix -c /opt/thermomix-master claude
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
