# webapp/deploy/thermomix-worker.service
[Unit]
Description=thermomix-master pipeline worker (processes pinned URLs)
After=thermomix-webapp.service

[Service]
Type=simple
WorkingDirectory=/opt/thermomix-master/webapp
Environment=NODE_ENV=production
Environment=THERMOMIX_STATE_DIR=/var/lib/thermomix-master
Environment=THERMOMIX_PYTHON=/opt/thermomix-master/.venv/bin/python
ExecStart=/usr/bin/node /opt/thermomix-master/webapp/worker/worker.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
